# builds/any/rootfs/jessie/common/overlay/etc/systemd/system/onl-bootd.service
[Unit]
Description=ONL boot.d

[Service]
Type=oneshot
ExecStart=/etc/boot.d/boot

[Install]
WantedBy=multi-user.target
